# nethunter-pro/overlays/husky/etc/systemd/system/firmware-modem.mount
[Unit]
Description=FUSE mount of Shannon 5300 modem firmware partition (/firmware/modem)
Documentation=https://github.com/mikethi/nhpro-native-husky
# Source: init.zuma.rc  "mount_all /vendor/etc/fstab.modem --early"
#         cbd.rc        "vendor.cbd.partition modem${slot_suffix}"
#
# The Samsung Shannon 5300 modem firmware is stored on the modem_a/modem_b
# UFS partition (MMIO 0x13200000, platform device 13200000.ufs).
# The nhpro-modem-mount helper detects the active A/B slot from /proc/cmdline
# and FUSE-mounts the correct ext4 partition at /firmware/modem using fuse2fs.
#
# ModemManager and the cpif/modem_if kernel driver read modem firmware from
# /firmware/modem after this unit is active.
After=systemd-udev-settle.service
Before=ModemManager.service
Wants=ModemManager.service

[Service]
Type=forking
ExecStartPre=/bin/mkdir -p /firmware/modem
ExecStart=/usr/local/bin/nhpro-modem-mount
ExecStop=/bin/umount /firmware/modem
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
